# prologin/.config/autostart/prologin.desktop
[Desktop Entry]
Exec=/usr/bin/firefox
Name=Prologin
Type=Application
Version=1.0
